# services/ffmpeg-stream.service
[Unit]
Description=FFmpeg Stream Service
After=network.target sound.target

[Service]
ExecStart=/usr/bin/ffmpeg -f alsa -i hw:0,0 -c:a libmp3lame -b:a 64k -f mp3 icecast://source:hackme@localhost:8000/stream
Restart=always
RestartSec=5
User=root
Group=audio
Environment="ALSA_CARD=default"

[Install]
WantedBy=multi-user.target